# files/etc/systemd/system/zero-key-policy.service
[Unit]
Description=Cardputer Zero internal keyboard policy
Documentation=file:/usr/local/bin/zero-key-policy
After=systemd-logind.service

[Service]
Type=simple
Environment=ZERO_KEY_POLICY_SEAT=seat-cardputer-zero
Environment=ZERO_KEY_POLICY_TTY=
ExecStart=/usr/local/bin/zero-key-policy
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
