# files/etc/systemd/system/zero-hdmi-lightdm-policy.service
[Unit]
Description=Cardputer Zero HDMI LightDM policy
After=systemd-udevd.service systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/zero-hdmi-lightdm-policy
TimeoutStartSec=25

[Install]
WantedBy=multi-user.target
