# files/etc/systemd/system/plymouth-quit-wait.service.d/10-cardputer-zero-order.conf
[Unit]
# Raspberry Pi OS may start plymouth-quit-wait.service before
# plymouth-quit.service. The wait helper then blocks forever and prevents
# multi-user.target/graphical.target from being reached. If the wait unit is in
# the boot transaction, make the quit unit part of the same transaction and run
# it first.
Requires=plymouth-quit.service
After=plymouth-quit.service
