# files/etc/systemd/system/lightdm.service.d/10-cardputer-zero-hdmi.conf
[Unit]
ConditionPathExists=/run/cardputer-zero/hdmi-connected
