# files/etc/systemd/system/lightdm.service.d/10-cardputer-zero-hdmi.conf
[Unit]
After=zero-hdmi-lightdm-policy.service
Wants=zero-hdmi-lightdm-policy.service
ConditionPathExists=/run/cardputer-zero/hdmi-connected
